# aosvc_win/awsao.service
[Unit]
Description=AWS AlwasyOn client service
After=dbus.service

[Service]
StandardError=journal
ExecStart=/usr/local/bin/awsao

[Install]
WantedBy=default.target